# source/debian/postinst
#!/bin/sh
# postinst script for logdata-anomaly-miner
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    analysisUser="aminer"
    analysisGroup="aminer"
# Prohibit read access to configuration for other processes
    if ! dpkg-statoverride --list /etc/aminer > /dev/null; then
      chown "root.${analysisGroup}" -- /etc/aminer
      chmod 00750 -- /etc/aminer
    fi
    if ! dpkg-statoverride --list /var/lib/aminer > /dev/null; then
      chmod 00700 -- /var/lib/aminer
      chown "${analysisUser}.${analysisGroup}" -- /var/lib/aminer
    fi
    ;;
esac

# We do not use pyc-files anyway, so disable generation.
alias pycompile=/bin/true

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# Start the service only when it was already enabled before updating.
if deb-systemd-helper debian-installed aminer.service && [ -d /run/systemd/system ]; then
	systemctl --system daemon-reload >/dev/null || true
	deb-systemd-invoke start aminer.service >/dev/null || true
fi
exit 0
