# source/debian/prerm
#!/bin/sh
# prerm script for logdata-anomaly-miner
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove)
# Stop all services and delete the users before removing the service.
# As quite some users won't use the daemon, ignore errors due
# to non-running service.
    if test -d /usr/lib/upstart; then
# Just stop, configuration file will be removed when package is
# removed.
        stop aminer || true
    elif test -e /bin/systemctl; then
        systemctl --quiet is-active aminer && systemctl stop aminer
        systemctl --quiet is-enabled aminer && systemctl disable aminer
    else
        invoke-rc.d aminer stop || true
    fi

# Revert overrides, otherwise next dpkg install will fail.
    dpkg-statoverride --remove /etc/aminer


# Delete user, will also delete group.
    userdel "aminer"
    ;;

    upgrade|deconfigure)
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
